# docker/supervisord.conf
[supervisord]
nodaemon=true
loglevel=info
logfile=/dev/null
logfile_maxbytes=0

[program:xvfb]
command=Xvfb :1 -screen 0 %(ENV_CORRIDORKEY_RESOLUTION)s -ac +extension GLX +render -noreset
autorestart=true
priority=10
startsecs=5
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0

[program:x11vnc]
command=x11vnc -display :1 -forever -shared -nopw -xkb -rfbport 5900
autorestart=true
priority=20
startsecs=2
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0

[program:novnc]
command=websockify --web=/usr/share/novnc 6080 localhost:5900
autorestart=true
priority=30
startsecs=3
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0

[program:filebrowser]
command=/usr/local/bin/filebrowser --port 6081 --address 0.0.0.0 --root /opt/corridorkey/ClipsForInference --noauth
autorestart=true
priority=35
startsecs=2
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0

[program:corridorkey]
command=/opt/corridorkey/.venv/bin/python /opt/corridorkey/main.py --gui
directory=/opt/corridorkey
environment=DISPLAY=":1",QT_QPA_PLATFORM="xcb",HF_HUB_DISABLE_PROGRESS_BARS="1"
autorestart=true
priority=40
startsecs=5
startretries=3
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0
